# systemd/config-server.service
[Unit]
Description=HiFiBerry Configuration Server
Documentation=man:config-server(1)
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
ExecStartPre=-/sbin/modprobe i2c-dev
ExecStartPre=/usr/bin/config-server --restore-settings
ExecStart=/usr/bin/config-server --host 0.0.0.0 --port 1081 --verbose
Restart=always
RestartSec=5
User=root
Group=root
PrivateMounts=no

# Capabilities needed for mount/unmount
CapabilityBoundingSet=CAP_SYS_ADMIN
DevicePolicy=auto

# Security
NoNewPrivileges=false
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=1024
LimitNPROC=64

[Install]
WantedBy=multi-user.target
