# systemd/config-server.service
[Unit]
Description=HiFiBerry Configuration Server
Documentation=man:config-server(1)
Wants=network-online.target config-detect.service
After=network-online.target systemd-modules-load.service config-detect.service create-uuid.service
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/bin/config-server --host 0.0.0.0 --port 1081 --verbose --auto-restore-settings
Restart=always
RestartSec=5
User=root
Group=audio
PrivateMounts=no

# Environment for PipeWire access
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="PIPEWIRE_RUNTIME_DIR=/run/pipewire"

# No CapabilityBoundingSet: config-server installs extensions via apt/dpkg, and
# package maintainer scripts need arbitrary root capabilities — apt's sandbox
# setuids to the _apt user (CAP_SETUID/SETGID), dpkg sets ownership and file
# capabilities (CAP_CHOWN/FOWNER/SETFCAP/DAC_OVERRIDE), and mounts still need
# CAP_SYS_ADMIN. Bounding to CAP_SYS_ADMIN made every apt operation fail
# ("seteuid 42 failed - Operation not permitted"). The service already runs as
# root with NoNewPrivileges=false, so an inclusive bounding set added negligible
# hardening in this posture. Leaving the directive out keeps the full root set.
DevicePolicy=auto

# Security - relaxed for PipeWire access
NoNewPrivileges=false
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectControlGroups=false

# Resource limits
LimitNOFILE=1024
LimitNPROC=64

[Install]
WantedBy=multi-user.target
