# systemd/config-server.service
[Unit]
Description=HiFiBerry Configuration Server
Documentation=man:config-server(1)
Wants=network-online.target config-detect.service
After=network-online.target systemd-modules-load.service config-detect.service create-uuid.service
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
ExecStart=/usr/bin/config-server --host 0.0.0.0 --port 1081 --verbose --auto-restore-settings
Restart=always
RestartSec=5
User=root
Group=audio
PrivateMounts=no

# Environment for PipeWire access
Environment="XDG_RUNTIME_DIR=/run/user/0"
Environment="PIPEWIRE_RUNTIME_DIR=/run/pipewire"

# Capabilities needed for mount/unmount
CapabilityBoundingSet=CAP_SYS_ADMIN
DevicePolicy=auto

# Security - relaxed for PipeWire access
NoNewPrivileges=false
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectControlGroups=false

# Resource limits
LimitNOFILE=1024
LimitNPROC=64

[Install]
WantedBy=multi-user.target
